# config/com.usetrack.menubar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.usetrack.menubar</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/yuan/bin/UseTrackMenuBar</string>
    </array>

    <!-- Start at login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart if crashes -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/tmp/usetrack-menubar.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/usetrack-menubar.stderr.log</string>

    <!-- GUI app needs access to the Aqua session -->
    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>
</dict>
</plist>
